# containers/php7/etc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:php]
command=php-fpm

[program:cron]
command=cron -f
